# .config/systemd/user/pipewire.service.d/override.conf
[Service]
CPUSchedulingPolicy=rr
CPUSchedulingPriority=88
